# systemd/smartmet-frontend.service
[Unit]

Description=Smartmet Weather Server Frontend

[Service]

# Not a real daemon
Type=simple

EnvironmentFile=/etc/smartmet/smartmet-frontend.defaults.env
EnvironmentFile=/etc/smartmet/smartmet-frontend.env

# Give capabilities what we need and restrict to them only
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

NoNewPrivileges=true

ExecStart=/usr/sbin/smartmetd --port=${PORT} --configfile ${CONFIGFILE}

LimitCORE=infinity
# LimitCORE=0

# Do not retry too fast if there are some I/O or similar issues
Restart=always
RestartSec=5s

# We do not signal completion to systemd
TimeoutStartSec=0

# Wait for 30s for the server to stop before sending SIGTERM and another 30s before SIGKILL
TimeoutStopSec=35s

# Disable swapping
MemorySwapMax=0M

LimitNOFILE=999999999

User=smartmet-server

# Run post step as root (prefix +) and target the main pid
ExecStartPost=+/bin/sh -eu -c "echo ${COREDUMP_FILTER} > /proc/$MAINPID/coredump_filter"

[Install]

WantedBy=multi-user.target
